# root/etc/uci-defaults/99_chongyoung
#!/bin/sh

# 配置迁移：添加新增的配置项
if ! uci -q get chongyoung.general.password_seed >/dev/null; then
    uci set chongyoung.general.password_seed=''
    uci commit chongyoung
fi

# 合并 opkg 备份的配置文件（如果存在）
if [ -f /etc/config/chongyoung-opkg ]; then
    # 临时重命名以便 uci 工具识别（uci 不支持带连字符的文件名作为配置名）
    mv /etc/config/chongyoung-opkg /etc/config/chongyoung_backup
    
    # 使用 uci 读取旧配置，完美支持多行值
    OLD_ENABLED=$(uci -q get chongyoung_backup.general.enabled)
    OLD_USERNAME=$(uci -q get chongyoung_backup.general.username)
    OLD_PWD_LIST=$(uci -q get chongyoung_backup.daily.password_list)
    
    # 恢复到新配置中
    [ -n "$OLD_ENABLED" ] && uci set chongyoung.general.enabled="$OLD_ENABLED"
    [ -n "$OLD_USERNAME" ] && uci set chongyoung.general.username="$OLD_USERNAME"
    [ -n "$OLD_PWD_LIST" ] && uci set chongyoung.daily.password_list="$OLD_PWD_LIST"
    
    uci commit chongyoung
    
    # 清理临时文件
    rm -f /etc/config/chongyoung_backup
fi

# 确保脚本有执行权限（虽然 opkg 会处理，但以防万一）
chmod +x /usr/bin/chongyoung.sh
chmod +x /usr/share/chongyoung/calc_pwd.lua

/etc/init.d/rpcd reload
/etc/init.d/chongyoung enable
/etc/init.d/chongyoung restart

exit 0
