# root/etc/uci-defaults/99_chongyoung
#!/bin/sh
/etc/init.d/rpcd reload
/etc/init.d/chongyoung enable
/etc/init.d/chongyoung restart
exit 0
